# 2root.smartd/etc/systemd/system/smartd.service.d/override.conf
[Unit]
After=user@1000.service
Wants=user@1000.service
