# 2root.ollama.whisker/etc/systemd/system/ollama.service.d/override.conf
[Service]
# Store models in /home (large btrfs partition) instead of /var/lib/ollama
Environment="OLLAMA_MODELS=/home/ollama/models"
# Bind to all interfaces so containers (podman) can reach it
Environment="OLLAMA_HOST=0.0.0.0:11434"
# Reduced context window for laptop GPU
Environment="OLLAMA_CONTEXT_LENGTH=8192"
# Allow /home to be accessible (default service has ProtectHome=yes)
ProtectHome=no
