# borg/.config/systemd/user/borg.timer
[Unit]
Description=Run borg backup daily.

[Timer]
Persistent=true
OnCalendar=daily
Unit=borg.service

[Install]
WantedBy=timers.target
