# 2target=root/smartd_simple_scrubbing/etc/systemd/system/raid_scrubbing.timer
[Unit]
Description=Run recollindex daily and every 3 hours

[Timer]
Unit=raid_scrubbing.service
Persistent=true
OnCalendar=daily
RandomizedDelaySec=3min

[Install]
WantedBy=timers.target
