# 2root.ollama/etc/systemd/system/ollama.service.d/override.conf
[Unit]
Requires=home.mount
After=home.mount

[Service]
# Store models in /home (large btrfs partition) instead of /var/lib/ollama
Environment="OLLAMA_MODELS=/home/ollama/models"
# Bind to all interfaces so containers (podman) can reach it
Environment="OLLAMA_HOST=0.0.0.0:11434"
# Extended context window — A4500 24GB VRAM
Environment="OLLAMA_CONTEXT_LENGTH=65536"
# Allow /home to be accessible (default service has ProtectHome=yes)
ProtectHome=no
