# recoll/.config/systemd/user/recollindex.service
[Unit]
Description=recollindex

[Service]
Type=oneshot
ExecStart=/usr/bin/recollindex > /tmp/recoll.out 2>&1
