# recoll/.config/systemd/user/recollindex.service
[Unit]
Description=recollindex

[Service]
;Type=oneshot
ExecStart=/home/dan/.local/bin/recollindex.sh
