# 2root.systemd/etc/systemd/system/failure-notification@.service.d/override.conf
[Unit]
Description=Send a notification about a failed systemd unit %i
After=network-online.target
Wants=network-online.target
