# 2root.raid_mdadm_scrubbing/etc/systemd/system/raid_scrubbing.service
[Unit]
Description=Simple RAID scrubbing
OnFailure=status-email-user@%n.service

[Service]
;Type=oneshot
ExecStart=/usr/local/bin/raid_scrubbing.bash
